# tt.conf
[program:tt]
command=python3 /home/curator2/tt/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/tt/test.err.log
stdout_logfile=/home/curator2/tt/test.out.log